# cloud-monitoring/config/cmon.conf
description     "C'mon"
author          "Darryl Ring <ringda@uvic.ca>"

start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/bin/python /usr/local/bin/cmon
