# docker/external_cloud/csv2/separated_containers/ansible_setup/cloud_scheduler/supervisord.conf
[supervisord]
nodaemon=true
 
[program:sshd]
command=/usr/sbin/sshd -p 947 -D
 
#[program:httpd]
#command=/bin/bash -c "rm -rf /run/httpd/* /tmp/httpd* && systemctl start httpd -DFOREGROUND"